# logit.service
[Unit]
Description=Logit service

[Service]
ExecStart=~/logit.sh

[Install]
WantedBy=multi-user.target
